# co2-zabbix.timer
[Unit]
Description=Timer fetch CO2 meter to zabbix server

[Timer]
OnUnitActiveSec=60seconds
OnBootSec=10seconds

[Install]
WantedBy=timers.target